# code/ch06/greeter/greeter.timer
# greeter.timer
[Unit]
Description=Runs Greeting service at the top of the hour

[Timer]
OnCalendar=hourly
