# code/ch06/greeter/greeter.service
# greeter.service
[Unit]
Description=My Greeting Service

[Service]
Type=oneshot
ExecStart=/home/mh9/greeting.sh

